# case-calendar.service.example
# systemd unit template for `case-calendar serve` (webhook receiver).
#
# Copy this file to /etc/systemd/system/case-calendar.service (or place it
# under ~/.config/systemd/user/ for a rootless setup), edit the User, paths,
# and any flags below, then:
#
#   systemctl daemon-reload
#   systemctl enable --now case-calendar.service
#   systemctl status case-calendar.service
#   journalctl -u case-calendar.service -f
#
# Assumed layout:
#   /opt/case-calendar/                  # the project (this repo, cloned)
#   /opt/case-calendar/.env              # COURTLISTENER_TOKEN, *_API_KEY,
#                                        # CASE_CALENDAR_WEBHOOK_SECRET, etc.
#   /opt/case-calendar/config.yaml       # cases / calendars / dockets
#   /opt/case-calendar/data/             # SQLite store (writable)
#   /opt/case-calendar/out/              # ICS + index.html (writable)
#   /opt/case-calendar/.case-calendar/   # OAuth token caches (writable)
#
# A dedicated unprivileged user owns the install:
#
#   useradd --system --home /opt/case-calendar --shell /usr/sbin/nologin casecalendar
#   chown -R casecalendar:casecalendar /opt/case-calendar
#
# Caddy fronts this on TCP/443 and reverse-proxies the webhook endpoint to
# 127.0.0.1:8000 (see the commented `webhook.example.com` block in
# Caddyfile.example). The serve process binds to localhost only by default,
# so it is never directly exposed to the public internet.
#
# OAuth setup (`case-calendar setup gcal` / `setup m365`) needs an
# interactive browser and CANNOT run under this unit. Either:
#   (a) run it as the casecalendar user on a machine with a browser, then
#       copy ~/.case-calendar/*.json into /opt/case-calendar/.case-calendar/
#   (b) SSH in with X11 forwarding (-X) and run `setup` once

[Unit]
Description=case-calendar webhook receiver
Documentation=https://github.com/seanpwhalen/case-calendar
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=casecalendar
Group=casecalendar
WorkingDirectory=/opt/case-calendar

# Loads COURTLISTENER_TOKEN, *_API_KEY, CASE_CALENDAR_WEBHOOK_SECRET, etc.
# Mode 0600, owned by the casecalendar user.
EnvironmentFile=/opt/case-calendar/.env

# OAuth token caches default to ~/.case-calendar/. Setting HOME here lets
# the defaults resolve to a path inside the install dir without needing
# config.yaml overrides. ProtectHome=true (below) blocks /home, /root, and
# /run/user — none of which we need with HOME pointed inside /opt.
Environment=HOME=/opt/case-calendar

# Adjust the uv path if `which uv` returns something different on your
# host (e.g. /root/.local/bin/uv if you installed uv as root, or
# /usr/bin/uv on distros that package it).
ExecStart=/usr/local/bin/uv run case-calendar serve --host 127.0.0.1 --port 8000

Restart=on-failure
RestartSec=5s

# --- Hardening ---
# Bare minimum — extend as paranoia permits. None of these break the
# expected runtime; they just shrink the unit's blast radius if something
# in the dependency tree gets compromised.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true

# The serve process writes to data/, out/, and .case-calendar/ — every
# other path stays read-only thanks to ProtectSystem=full.
ReadWritePaths=/opt/case-calendar/data
ReadWritePaths=/opt/case-calendar/out
ReadWritePaths=/opt/case-calendar/.case-calendar

[Install]
WantedBy=multi-user.target
